# crontab
* * * * * python3 /home/kontorshund/web/hello.py
# DON'T REMOVE LAST LINE